# deploy/macos/com.predigy.wx-stat-curate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Forecast-driven curator for Kalshi temperature markets.
         **Phase 2 NBM probabilistic path is on by default.** Pulls
         NBM 21-quantile model data from NOAA's S3 bucket, samples
         at each airport's grid cell, interpolates the CDF at the
         Kalshi-side threshold to get a real calibrated probability,
         and emits a StatRule per market.

         Output file (~/.config/predigy/wx-stat-rules.json) is
         consumed directly by the consolidated engine's wx-stat
         strategy when PREDIGY_WX_STAT_RULE_FILE points at it.
         Same-day/past temperature markets are gated through ASOS
         observed extremes before forecast/NBM scoring.

         Cadence: every 3h. NBM qmd publishes only at 00/06/12/18
         UTC (4 cycles/day), but with the 13h lookback the curator
         lands in a fully-published cycle every tick. The 3h cron
         schedule means every cycle gets re-run twice as fresh
         observations / Kalshi quotes come in.

         Pre-flight: the script refuses to start if NWS_USER_AGENT
         isn't set in ~/.zprofile (NOAA's bucket doesn't gate on
         it but it's polite + best practice).
    -->
    <key>Label</key>
    <string>com.predigy.wx-stat-curate</string>

    <key>Disabled</key>
    <false/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec /Users/dan/code/predigy/deploy/scripts/wx-stat-curate.sh</string>
    </array>

    <!-- Every 3h: 03/06/09/12/15/18/21/00 local. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <!-- Re-curate immediately on enable so the rule file is fresh
         the moment the operator flips Disabled=false. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/predigy/wx-stat-curate.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/predigy/wx-stat-curate.stderr.log</string>
</dict>
</plist>
