# deploy/macos/com.predigy.calibration.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Calibration evidence tick. Backfills public settlement outcomes,
         reconciles settled venue-flat stale DB positions, and writes
         calibration_reports. Does not enable rules or submit orders. -->
    <key>Label</key>
    <string>com.predigy.calibration</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec /Users/dan/code/predigy/deploy/scripts/predigy-calibration-run.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/predigy/calibration.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/predigy/calibration.stderr.log</string>
</dict>
</plist>
