# deploy/macos/com.predigy.calibration.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Calibration evidence tick. Backfills public settlement outcomes,
         reconciles settled venue-flat stale DB positions, and writes
         calibration_reports. Does not enable rules or submit orders. -->
    <key>Label</key>
    <string>com.predigy.calibration</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec /Users/dan/code/predigy/deploy/scripts/predigy-calibration-run.sh</string>
    </array>

    <!-- 2026-05-09: cadence 1h → 15m. The bundled `reconcile-venue-flat
         --write` step is the auto-resolver for phantom DB rows that the
         engine's internal reconciler observes-only. At 1h cadence the
         engine spams "position_mismatches" drift warnings for up to 60
         min after each game settles. 15 min cuts that to ≤15 min. The
         sync-settlements + report steps are wider work but light enough
         to also run 4x/hour without strain. -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/predigy/calibration.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/predigy/calibration.stderr.log</string>
</dict>
</plist>
