# deploy/macos/com.predigy.wx-stat-curate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Forecast-driven curator for Kalshi temperature markets.
         **Phase 2 NBM probabilistic path is on by default.** Pulls
         NBM 21-quantile model data from NOAA's S3 bucket, samples
         at each airport's grid cell, interpolates the CDF at the
         Kalshi-side threshold to get a real calibrated probability,
         and emits a StatRule per market.

         Output file (~/.config/predigy/wx-stat-rules.json) is
         consumed directly by the consolidated engine's wx-stat
         strategy when PREDIGY_WX_STAT_RULE_FILE points at it.
         Same-day/past temperature markets are gated through ASOS
         observed extremes before forecast/NBM scoring.

         Cadence: hourly. NBM qmd publishes only at 00/06/12/18
         UTC (4 cycles/day), but the curator's 13h lookback lands
         in a fully-published cycle every tick. Hourly means each
         NBM cycle gets re-run ~6 times as observed-extreme gating
         and Kalshi quotes refresh; same-day temperature markets
         that hit observation extremes mid-day get re-evaluated
         within an hour rather than waiting for the next 3h cron.
         A no-change tick is cheap — pure REST + S3 + DB writes.

         Pre-flight: the script refuses to start if NWS_USER_AGENT
         isn't set in ~/.zprofile (NOAA's bucket doesn't gate on
         it but it's polite + best practice).
    -->
    <key>Label</key>
    <string>com.predigy.wx-stat-curate</string>

    <key>Disabled</key>
    <false/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec /Users/dan/code/predigy/deploy/scripts/wx-stat-curate.sh</string>
    </array>

    <!-- Hourly. StartInterval rather than StartCalendarInterval
         keeps it simple. Combined with RunAtLoad below, the rule
         file is never more than ~1h stale (and is fresh
         immediately after a re-install / boot). -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Re-curate immediately on enable so the rule file is fresh
         the moment the operator flips Disabled=false. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/predigy/wx-stat-curate.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/predigy/wx-stat-curate.stderr.log</string>
</dict>
</plist>
