# deploy/linux/systemd/predigy-db-backup.service
[Unit]
Description=predigy daily pg_dump backup with rotation
After=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=%h/code/predigy
EnvironmentFile=%h/.config/predigy/env
ExecStart=%h/code/predigy/deploy/scripts/db-backup.sh
StandardOutput=append:%h/.local/state/predigy/logs/db-backup.stdout.log
StandardError=append:%h/.local/state/predigy/logs/db-backup.stderr.log
