# deploy/linux/systemd/predigy-engine.service
[Unit]
# Consolidated predigy trading engine. Direct counterpart of
# deploy/macos/com.predigy.engine.plist on macOS. KeepAlive replaced
# by Restart=on-failure; the launchd ThrottleInterval=30s
# translates to RestartSec=30. The engine refuses to start without
# its env vars (KALSHI_KEY_ID, KALSHI_PEM), which load from
# %h/.config/predigy/env via EnvironmentFile=.
Description=predigy trading engine (consolidated)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/code/predigy
EnvironmentFile=%h/.config/predigy/env
ExecStart=%h/code/predigy/deploy/scripts/engine-run.sh
Restart=on-failure
RestartSec=30
# Cap restart loops the way launchd's ThrottleInterval did.
StartLimitIntervalSec=600
StartLimitBurst=5
# Stream stderr to a tail-able file (the dashboard reads strategy
# logs via file path; same file scheme as macOS, just under XDG).
StandardOutput=append:%h/.local/state/predigy/logs/engine.stdout.log
StandardError=append:%h/.local/state/predigy/logs/engine.stderr.log

[Install]
WantedBy=default.target
